# etc/services.d/elasticsearch/run
#!/usr/bin/with-contenv /bin/ash
exec s6-setuidgid elasticsearch /opt/containerpilot/containerpilot /opt/elasticsearch/bin/elasticsearch
